# supervisor/jahshaka-api.ini
[program:jahshaka-auth-server]
directory=/var/www/aspnet/Jahshaka.API
command=dotnet Jahshaka.API.dll
process_name=%(program_name)s
autostart=true
autorestart=true
stderr_logfile=/var/log/jahshaka-api.err.log
stdout_logfile=/var/log/jahshaka-api.out.log
environment=ASPNETCORE_ENVIRONMENT=Staging
user=root
stopsignal=INT
